# config/systemd/magic-workflow.service
# Magic Workflow — run the full stack on boot.
# Installed and path-adjusted by scripts/install-server.sh, or do it manually:
#   sudo cp config/systemd/magic-workflow.service /etc/systemd/system/
#   sudo sed -i "s|/opt/magic-workflow|$(pwd)|" /etc/systemd/system/magic-workflow.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now magic-workflow
#
# Podman hosts: use the RHEL path (rootful podman + this unit works if you change
# the ExecStart to `podman-compose`, or use `podman generate systemd` quadlets).
[Unit]
Description=Magic Workflow collaboration suite
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/magic-workflow
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
